# .nova_sync/com.nova.pulse-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nova.pulse-monitor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/liman/Desktop/huntlink-stable-final/.nova_sync/nova-pulse-monitor.sh</string>
        <string>start</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/liman/Desktop/huntlink-stable-final/.nova_sync</string>
    
    <key>StandardOutPath</key>
    <string>/Users/liman/Desktop/huntlink-stable-final/.nova_sync/monitor.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/liman/Desktop/huntlink-stable-final/.nova_sync/monitor-error.log</string>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>